# etc/udev/rules.d/Wacom-Intuos-BT-S.rules
ACTION=="add" ATTRS{idVendor}=="056a", ATTRS{idProduct}=="03c5", RUN+="/usr/bin/wacom-config"
